# deploy/systemd/algua-paper.service
[Unit]
Description=Algua paper-operator cycle (session-gated)
# Fired by algua-paper.timer. The XNYS calendar gate + per-session idempotency marker inside
# `algua operator run` decide whether the driver below actually runs; a weekend/holiday firing or a
# re-fire of an already-completed session no-ops.

[Service]
Type=oneshot
EnvironmentFile=/etc/algua/algua.env
WorkingDirectory=/opt/algua
# Backstop for a wedged run: `operator run` already kills a hung driver at its own wall-clock cap
# (2x the paper job's 900s grace = 1800s) and releases operator.lock, but if the wrapper ITSELF
# hangs (e.g. blocked on the lock or a stuck alert hook), systemd SIGTERMs the unit here so a firing
# can never wedge forever. Set ABOVE the app-level cap so the app timeout (which alerts + leaves the
# marker unwritten for a clean retry) fires first; this is only the last resort.
TimeoutStartSec=2400
# The command after `--` MUST exactly match the `paper` job's canonical argv template
# (`algua paper run-all --snapshot {snapshot}`) or `operator run` fail-closes with command_mismatch.
# ${ALGUA_PAPER_SNAPSHOT} is expanded by systemd from the EnvironmentFile. `algua` must be on PATH
# for the operator's subprocess (a `.venv/bin` entry or a symlink); do NOT append ad-hoc flags — a
# drawdown/window override is a manual human path OUTSIDE the timer (see README, direct-invocation).
ExecStart=/opt/algua/.venv/bin/algua operator run --job paper -- algua paper run-all --snapshot ${ALGUA_PAPER_SNAPSHOT}
